# 2.x/config/supervisor/crond.ini
[program:crond]
command=/usr/sbin/crond -n -p
redirect_stderr=true