# 2.x/config/supervisor/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
redirect_stderr=true

[program:php-fpm]
command=/usr/sbin/php-fpm -F
redirect_stderr=true