# 1.x/config/supervisor/cronpy.ini
[program:cronpy]
command=/usr/local/bin/cronpy
redirect_stderr=true