# 1.x/config/supervisor/httpd.ini
[program:httpd]
command=/usr/sbin/httpd -T -DFOREGROUND
redirect_stderr=true